# debian/sitter-backup-etc.cron.weekly
#!/bin/sh -e
#
# Run a commit of the latest changes

cd /etc
git add .
git commit --quiet --message="/etc auto-commit" .

# vim: ts=4 sw=4 et
